# daemon-py/systemd/valueguard.service
# ValueGuard daemon — per-user systemd unit (Linux).
#
# Install:
#   mkdir -p ~/.config/systemd/user
#   cp valueguard.service ~/.config/systemd/user/
#   systemctl --user daemon-reload && systemctl --user enable --now valueguard
#
# Runs in the user session (never a system service — it must see the user's
# display). DISPLAY/XAUTHORITY must match the session's X server; on a
# desktop login they're usually in the systemd user environment already
# (`systemctl --user show-environment`), otherwise import them:
#   systemctl --user import-environment DISPLAY XAUTHORITY

[Unit]
Description=ValueGuard on-device content filter (log-only)
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=exec
# Adjust paths: the venv with daemon-py/requirements.txt installed, and the
# exported model (model-conversion/output/SigLIP2Vision.fp32.onnx).
ExecStart=%h/valueguard/.venv/bin/python -m valueguard_daemon --model %h/valueguard/models/SigLIP2Vision.fp32.onnx
WorkingDirectory=%h/valueguard
Restart=on-failure
RestartSec=5
Nice=10

[Install]
WantedBy=graphical-session.target
